# pi-client/ahd-pi.service
[Unit]
Description=AHD Radio DJ Pi Client
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/default/ahd-pi
WorkingDirectory=/opt/ahd-pi-client/pi-client
ExecStart=/usr/bin/node client.js
Restart=always
RestartSec=3
User=pi

[Install]
WantedBy=multi-user.target
